# koel/m3usync/m3usync-supervisor.conf
[program:m3usync]
user = root
command = /root/m3usync.sh
umask = 000

